# ops/deploy/research-os-a1-final.service
[Unit]
Description=Dual-Contour Research OS final offline A1 Bridge
Documentation=https://github.com/KirPonomarev/dual-contour-research-os
After=docker.service
Requires=docker.service
ConditionPathIsSocket=%t/docker.sock

[Service]
Type=simple
Environment=DOCKER_HOST=unix://%t/docker.sock
Environment=RESEARCH_OS_ENVIRONMENT=pre-soak
Environment=RESEARCH_OS_EXTERNAL_ACTION_AUTHORITY=false
ExecStartPre=/bin/sh -eu -c '! systemctl --user is-active --quiet research-os-bridge.service; ! systemctl --user is-enabled --quiet research-os-bridge.service; state=$(/usr/bin/docker container inspect research-os-bridge --format "{{json .State.Running}}" 2>/dev/null || true); test "$state" != true'
ExecStartPre=/usr/bin/docker run --rm --network=none --entrypoint=/bin/sh --user=0:0 --mount=type=volume,source=research-os-a1-runtime,target=/var/lib/research-os @@IMAGE_ID@@ -eu -c 'chown 10001:10001 /var/lib/research-os; chmod 0710 /var/lib/research-os; install -m 0600 -o 10001 -g 10001 /dev/null /var/lib/research-os/.runtime-initialized'
ExecStartPre=-/usr/bin/docker container inspect research-os-a1-bridge
ExecStartPre=-/usr/bin/docker container create --name=research-os-a1-bridge --user=10001:10001 --group-add=10001 --network=none --read-only --cap-drop=ALL --security-opt=no-new-privileges:true --pids-limit=256 --memory=2147483648 --cpus=2 --restart=no --env=RESEARCH_OS_ENVIRONMENT=pre-soak --env=RESEARCH_OS_EXTERNAL_ACTION_AUTHORITY=false --label=org.research-os.release-sha=@@RELEASE_SHA@@ --label=org.research-os.policy-sha256=@@POLICY_SHA256@@ --label=org.research-os.config-sha256=@@CONFIG_SHA256@@ --label=org.research-os.principal-role-map=10001:operator,10002:collector,10003:scout,10004:connected_worker --label=org.research-os.control-socket-owner=10001:10001 --label=org.research-os.control-socket-mode=0660 --label=org.research-os.runtime-parent-mode=0710 --label=org.research-os.control-transport=AF_UNIX --mount=type=volume,source=research-os-a1-runtime,target=/var/lib/research-os --mount=type=volume,source=research-os-a1-config,target=/run/research-os,readonly @@IMAGE_ID@@
ExecStartPre=-/usr/bin/docker stop --time=30 research-os-a1-bridge
ExecStart=/usr/bin/docker start --attach research-os-a1-bridge
ExecStop=/usr/bin/docker stop --time=30 research-os-a1-bridge
Restart=on-failure
RestartSec=5s
TimeoutStartSec=90s
TimeoutStopSec=45s
NoNewPrivileges=yes
RestrictAddressFamilies=AF_UNIX
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
TasksMax=256
MemoryMax=2147483648
CPUQuota=200%
UMask=0077

[Install]
WantedBy=default.target
