# ops/deploy/research-os-bridge.service
[Unit]
Description=Dual-Contour Research OS offline Bridge
Documentation=https://github.com/KirPonomarev/dual-contour-research-os
After=docker.service
Requires=docker.service
ConditionPathIsSocket=%t/docker.sock

[Service]
Type=simple
Environment=DOCKER_HOST=unix://%t/docker.sock
Environment=RESEARCH_OS_ENVIRONMENT=pre-soak
Environment=RESEARCH_OS_EXTERNAL_ACTION_AUTHORITY=false
ExecStartPre=-/usr/bin/docker container inspect research-os-bridge
ExecStartPre=-/usr/bin/docker container create --name=research-os-bridge --user=10001:10001 --network=none --read-only --cap-drop=ALL --security-opt=no-new-privileges:true --pids-limit=256 --memory=2147483648 --cpus=2 --restart=unless-stopped --env=RESEARCH_OS_ENVIRONMENT=pre-soak --env=RESEARCH_OS_EXTERNAL_ACTION_AUTHORITY=false --label=org.research-os.release-sha=@@RELEASE_SHA@@ --label=org.research-os.policy-sha256=@@POLICY_SHA256@@ --label=org.research-os.config-sha256=@@CONFIG_SHA256@@ --mount=type=volume,source=research-os-bridge-runtime,target=/var/lib/research-os --mount=type=volume,source=research-os-bridge-config,target=/run/research-os,readonly @@IMAGE_ID@@
ExecStartPre=-/usr/bin/docker stop --time=30 research-os-bridge
ExecStart=/usr/bin/docker start --attach research-os-bridge
ExecStop=/usr/bin/docker stop --time=30 research-os-bridge
Restart=on-failure
RestartSec=5s
TimeoutStartSec=90s
TimeoutStopSec=45s
NoNewPrivileges=yes
RestrictAddressFamilies=AF_UNIX
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateDevices=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
TasksMax=256
MemoryMax=2147483648
CPUQuota=200%
UMask=0077

[Install]
WantedBy=default.target
